# ansible/roles/fuse/templates/fuse.service.j2
[Unit]
Description= fune service
After=network.target

[Service]

Type=oneshot

ExecStart=/opt/fuse/mount_fuse.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
