# ansible/roles/jupyterhub/templates/jupyterhub.service.j2
[Unit]
Description=Jupyterhub service
After=network.target

[Service]
ExecStart=jupyterhub -f /etc/jupyterhub/jupyterhub_config_simple.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
